# pbbacklight.service.in
[Unit]
Description=Pixelbook backlight driver
ConditionPathExists=/sys/class/drm_dp_aux_dev
ConditionPathExists=/sys/class/backlight/intel_backlight

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/pbbacklight
Restart=always

[Install]
WantedBy=multi-user.target
