# pbkbd.service.in
[Unit]
Description=Pixelbook Keyboard Driver
Wants=pbkbd-backlight.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/pbkbd

[Install]
WantedBy=getty.target
WantedBy=graphical.target
