# pbkbd-backlight.service
[Unit]
Description=Pixelbook Keyboard Backlight Driver

[Service]
ExecStart=/usr/local/sbin/pbkbd-backlight

[Install]
WantedBy=pbkbd.service
