# pbkbd-backlight.service.in
[Unit]
Description=Pixelbook Keyboard Backlight Driver
ConditionPathExists=/sys/class/leds/chromeos::kbd_backlight
ConditionPathExists=/sys/bus/iio/devices

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/pbkbd-backlight
