# pbkbd.service.in
[Unit]
Description=Pixelbook Keyboard Driver

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/pbkbd

[Install]
WantedBy=getty.target
WantedBy=graphical.target
Also=pbkbd-backlight.service
